# systemd/openclaw-gateway.service
[Unit]
Description=Openclaw Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=UNDICI_NO_HTTP2=1
EnvironmentFile=-%h/.config/openclaw/env
ExecStart=/usr/local/bin/openclaw gateway start
Restart=on-failure
RestartSec=5
WorkingDirectory=%h

[Install]
WantedBy=default.target
